# cron.d/dmponline.cron
SHELL=/bin/bash
BASH_ENV=/opt/DMPonline_v4/env.sh
MAILTO=nicolas.franck@ugent.be
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/rvm/bin

0 20 * * * dmponline /opt/DMPonline_v4/bin/daily_export_projects
0 0 * * * dmponline /opt/DMPonline_v4/bin/cleanup_locks
