# cron.d/dmponline.cron
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/rvm/bin

0 0 * * * dmponline RAILS_ENV=production /opt/DMPonline_v4/bin/daily_export /opt/dmponline_internal 2> /var/log/dmponline/daily_export.error
